# root/etc/services.d/ngircd/run
#!/usr/bin/with-contenv bash


exec \
	/usr/sbin/ngircd -n -f /config/ngircd.conf
